# etc/cron.d/sprinklerswitch
0 */12 * * * pi /bin/node /opt/sprinklerswitch/sprinkler.js /etc/sprinklerswitch/config.json 2>&1 | /usr/bin/logger -t SPRINKLERSWITCH
