# etc/cron.d/sprinkler
[ -r /etc/default/$NAME ] && . /etc/default/$NAME

0 23/2 * * * pi $NODE_HOME/bin/node /opt/sprinkler/sprinkler.js "$LATITUDE" "$LONGITUDE" 2>&1 | /usr/bin/logger -t SPRINKLER
